# templates/etc/profile.d/pipx.sh.j2
# {{ ansible_managed }}
export PIPX_BIN_DIR="/usr/local/bin"
export PIPX_HOME="/usr/local/share/pipx"
